# google_wipe.service
[Unit]
Description=chrome wipe script

[Service]
ExecStart=/usr/local/bin/google_wipe.sh

[Install]
WantedBy=multi-user.target
